# cron.daily/speedtest-register
#!/bin/sh
/usr/bin/php /var/www/speedtest/bin/console speedtest:register
